# agent/systemd/shaka-payment.service
[Unit]
Description=Shaka Payment Service (Stripe Terminal)
After=network.target shaka-vend.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/shaka
EnvironmentFile=/etc/default/shaka-stripe
ExecStart=/usr/bin/python3 /home/shaka/shaka_payment_service.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
