# agent/systemd/shaka-nayax.service
[Unit]
Description=Shaka Nayax Payment Service (Marshall Protocol)
After=network.target shaka-vend.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/shaka
EnvironmentFile=/etc/default/shaka-nayax
ExecStart=/usr/bin/python3 /home/shaka/shaka_nayax_service.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
