# agent/systemd/shaka-camera.service
[Unit]
Description=Shaka Camera Server
After=network.target

[Service]
Type=simple
User=shaka
WorkingDirectory=/home/shaka
ExecStart=/usr/bin/python3 /home/shaka/camera_server.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
